# deploy/pocket-tracker.service
[Unit]
Description=ptcgp-tracker web server
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=www-data
WorkingDirectory=/opt/ptcgp-tracker
Environment=NODE_ENV=production
EnvironmentFile=-/opt/ptcgp-tracker/.env
ExecStart=/usr/bin/npm run start
Restart=on-failure
RestartSec=5
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
